# setup/redis/systemd/system/redis.target
[Unit]
Description=Redis target allowing to start/stop all redis@.service instances at once
